# deploy/openclaw.service
[Unit]
Description=OpenClaw VRM Operator Gateway
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/openclaw
ExecStartPre=/opt/openclaw/deploy/fetch-secrets.sh jarvis-ml-dev
ExecStartPre=/opt/openclaw/deploy/generate-config.sh
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
